# gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn hello:app -c gunicorn.py
directory=/deploy/app
